# installer/build/scripts/systemd/engine_installer/engine_installer_startup.service
[Unit]
Description=Demo VCH Installer Wizard Startup Configuration
Documentation=https://github.com/vmware/vic
After=network-online.target systemd-resolved.service iptables.service ova-firewall.service fileserver.service data.mount
Before=engine_installer.service
Requires=iptables.service fileserver.service data.mount

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /etc/vmware/engine_installer/configure_engine_installer.sh

[Install]
WantedBy=multi-user.target
