# installer/build/scripts/systemd/engine_installer/engine_installer.service
[Unit]
Description=Demo VCH Installer Wizard Service
Documentation=https://github.com/vmware/vic
After=vic-appliance.target fileserver.service
Requires=vic-appliance.target fileserver.service

[Service]
Type=simple
Restart=on-failure
RestartSec=15
ExecStartPre=/usr/bin/bash /etc/vmware/engine_installer/configure_engine_installer.sh
ExecStart=/etc/vmware/engine_installer/start_engine_installer.sh

[Install]
WantedBy=vic-components.target
