# installer/build/scripts/systemd/engine_installer/engine_installer_startup.service
[Unit]
Description=Demo VCH Installer Wizard Startup Configuration
Documentation=https://github.com/vmware/vic
After=vic-appliance.target fileserver.service
Requires=vic-appliance.target fileserver.service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /etc/vmware/engine_installer/configure_engine_installer.sh

[Install]
WantedBy=multi-user.target
