# installer/build/scripts/systemd/psc/get_token.timer
[Unit]
Description=Run the PSC Get Token service every hour
Documentation=http://github.com/vmware/vic-product/installer

[Timer]
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
