# installer/build/scripts/systemd/psc/get_token.service
[Unit]
Description=PSC Get Token
Documentation=http://github.com/vmware/vic-product/installer
After=vic-appliance-ready.target get_token.path
Requires=vic-appliance-ready.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /etc/vmware/psc/get_token.sh

[Install]
WantedBy=multi-user.target
