# installer/build/scripts/systemd/psc/get_token.path
[Unit]
Description=Get PSC Tokens
Documentation=https://github.com/vmware/vic-product/installer

[Path]
PathExists=/etc/vmware/psc/admiral/psc-config.properties

[Path]
PathChanged=/etc/vmware/psc/admiral/psc-config.properties
Unit=harbor_startup.service

[Path]
PathChanged=/etc/vmware/psc/admiral/psc-config.properties
Unit=admiral_startup.service

[Install]
WantedBy=paths.target
