# installer/build/scripts/systemd/psc/get_token.service
[Unit]
Description=PSC Get Token
Documentation=http://github.com/vmware/vic-product/installer
Before=harbor.service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /etc/vmware/psc/get_token.sh

[Install]
WantedBy=multi-user.target
