# installer/build/scripts/systemd/units/configure-rsyslog.service
[Unit]
Description=Configure Remote Syslog Server
Before=rsyslog.service
Requires=vic-appliance-environment.service
After=vic-appliance-environment.service

[Service]
Type=oneshot
EnvironmentFile=/etc/vmware/environment
ExecStart=/etc/vmware/configure-rsyslog.sh

[Install]
WantedBy=rsyslog.service
