# installer/build/scripts/systemd/units/vic-appliance-rsyslog-localfiles.service
[Unit]
Description=Configure forwarding local files
After=admiral.service vic-machine-server.service
Requires=admiral.service vic-machine-server.service

[Service]
Type=oneshot
EnvironmentFile=/etc/vmware/environment
ExecStart=/etc/vmware/configure-rsyslog-localfiles.sh

[Install]
WantedBy=vic-appliance.target
