# installer/build/scripts/systemd/units/vic-appliance.target
[Unit]
Description=System customization for VIC Appliance
Documentation=https://github.com/vmware/vic-product/blob/master/installer/docs/DESIGN.md

Requires=multi-user.target
After=multi-user.target

AllowIsolate=yes
