# installer/build/scripts/systemd/units/vic-appliance-environment.service
[Unit]
Description=VIC Appliance Environment
Documentation=https://github.com/vmware/vic-product

[Service]
Type=oneshot
ExecStart=/etc/vmware/vic-appliance-environment.sh
Requires=ovf-network.service
After=ovf-network.service

[Install]
WantedBy=vic-appliance-ready.target
