# installer/build/scripts/harbor/harbor.service
[Unit]
Description=Harbor Container Registry
Documentation=https://vmware.github.io/vic-product/index.html#getting-started
After=vic-appliance-ready.target psc-ready.target admiral.service
Requires=vic-appliance-ready.target psc-ready.target

[Service]
Type=simple
Restart=on-failure
RestartSec=15
EnvironmentFile=/etc/vmware/environment
ExecStartPre=-/usr/local/bin/docker-compose -f /etc/vmware/harbor/docker-compose.yml -f /etc/vmware/harbor/docker-compose.notary.yml -f /etc/vmware/harbor/docker-compose.clair.yml down -v
ExecStartPre=-/usr/local/bin/docker-compose -f /etc/vmware/harbor/docker-compose.yml -f /etc/vmware/harbor/docker-compose.notary.yml -f /etc/vmware/harbor/docker-compose.clair.yml rm -f
ExecStartPre=-/usr/bin/docker stop dch-push
ExecStartPre=-/usr/bin/docker rm dch-push
ExecStartPre=/usr/bin/bash /etc/vmware/harbor/configure_harbor.sh
ExecStartPre=/usr/bin/python /etc/vmware/harbor/prepare --conf /storage/data/harbor/harbor.cfg --with-notary --with-clair
ExecStart=/etc/vmware/harbor/start_harbor.sh
ExecStartPost=/usr/bin/bash /etc/vmware/harbor/configure_internal_network.sh
ExecStop=/usr/local/bin/docker-compose -f /etc/vmware/harbor/docker-compose.yml -f /etc/vmware/harbor/docker-compose.notary.yml -f /etc/vmware/harbor/docker-compose.clair.yml down -v
ExecStopPost=/usr/local/bin/docker-compose -f /etc/vmware/harbor/docker-compose.yml -f /etc/vmware/harbor/docker-compose.notary.yml -f /etc/vmware/harbor/docker-compose.clair.yml rm -f

[Install]
WantedBy=vic-appliance.target
