# installer/packer/scripts/systemd/harbor/harbor_startup.path
[Unit]
Description=Harbor PSC token dependency
Documentation=https://github.com/vmware/vic-product/installer

[Path]
PathExists=/data/harbor/tokens.properties

[Install]
WantedBy=multi-user.target
