# installer/build/scripts/vic-machine-server/vic-machine-server.service
[Unit]
Description=VIC Machine API
Documentation=https://github.com/vmware/vic
After=vic-appliance-ready.target
Requires=vic-appliance-ready.target

[Service]
Type=simple
Restart=on-failure
RestartSec=15
EnvironmentFile=/etc/vmware/environment
ExecStartPre=-/usr/bin/docker kill vic-machine-server
ExecStartPre=-/usr/bin/docker rm -f vic-machine-server
ExecStartPre=/etc/vmware/vic-machine-server/configure-vic-machine-server.sh
ExecStart=/usr/bin/docker run --cap-drop ALL --cap-add NET_BIND_SERVICE --rm --name vic-machine-server -v /storage/data/certs:/certs:ro -v /storage/log/vic-machine-server:/var/log/vic-machine-server -p ${VIC_MACHINE_SERVER_PORT}:443 vmware/vic-machine-server:ova
ExecStop=/usr/bin/docker stop vic-machine-server

[Install]
WantedBy=vic-appliance.target
