# installer/packer/scripts/systemd/kov/kov_env_vars.service
[Unit]
Description=Kubernetes on Vsphere
Documentation=https://github.com/vmware/kov
After=network-online.target admiral_startup.service fileserver_startup.service admiral.service engine_installer_startup.service
Requires=docker.service admiral_startup.service fileserver_startup.service admiral.service engine_installer_startup.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/etc/vmware/kov/configure_kov_env_vars.sh

[Install]
WantedBy=multi-user.target
