# installer/build/scripts/systemd/fileserver/fileserver_startup.service
[Unit]
Description=VIC Unified Installer Web Server Startup Configuration
Documentation=https://github.com/vmware/vic
After=network-online.target systemd-resolved.service iptables.service ova-firewall.service network-online.target data.mount
Before=fileserver.service
Requires=iptables.service network-online.target data.mount

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /etc/vmware/fileserver/configure_fileserver.sh

[Install]
WantedBy=multi-user.target
