# installer/build/scripts/systemd/fileserver/fileserver_startup.service
[Unit]
Description=VIC Unified Installer Web Server Startup Configuration
Documentation=https://github.com/vmware/vic
After=vic-appliance.target
Requires=vic-appliance.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /etc/vmware/fileserver/configure_fileserver.sh

[Install]
WantedBy=multi-user.target
