# installer/build/scripts/systemd/fileserver/fileserver.service
[Unit]
Description=VIC Unified Installer Web Server
Documentation=https://github.com/vmware/vic
After=vic-appliance-ready.target fileserver_startup.service
Requires=vic-appliance-ready.target

[Service]
Type=simple
Restart=on-failure
RestartSec=5
ExecStartPre=-/usr/bin/systemctl stop landing_server.service
ExecStart=/etc/vmware/fileserver/start_fileserver.sh

[Install]
WantedBy=multi-user.target
