# installer/build/scripts/systemd/appliance/firstboot.path
[Unit]
Description=Docker images and chpasswd dependency
Documentation=https://github.com/vmware/vic-product/installer

[Path]
PathExists=/etc/vmware/firstboot

[Install]
WantedBy=multi-user.target
