# installer/build/scripts/systemd/appliance/firstboot.service
[Unit]
Description=Loads Cached Docker Images
After=data.mount docker.service
Wants=data.mount docker.service

[Service]
Type=oneshot
ExecStart=/etc/vmware/firstboot.sh

[Install]
WantedBy=vic-appliance.target
