# installer/build/scripts/systemd/appliance/firstboot.service
[Unit]
Description=Loads Cached Docker Images
After=vic-mounts.target docker.service
Wants=vic-mounts.target docker.service

[Service]
Type=oneshot
ExecStart=/etc/vmware/firstboot.sh
ExecStartPost=/usr/bin/systemctl disable firstboot.service

[Install]
WantedBy=vic-appliance.target
