# installer/build/scripts/systemd/appliance/vic-appliance-docker-images-loaded.path
[Unit]
Description=Docker images loaded
Documentation=https://github.com/vmware/vic-product/installer

[Path]
PathExists=/etc/vmware/firstboot
Unit=vic-appliance.target

[Install]
WantedBy=multi-user.target
