# installer/build/scripts/systemd/appliance/vic-appliance-load-docker-images.service
[Unit]
Description=Load Cached Docker Images
After=vic-mounts.target docker.service
Requires=vic-mounts.target docker.service

[Service]
Type=oneshot
ExecStart=/etc/vmware/load-docker-images.sh
ExecStartPost=/usr/bin/systemctl disable vic-appliance-load-docker-images.service

[Install]
WantedBy=vic-appliance.target
