# installer/build/scripts/systemd/admiral/admiral_startup.service
[Unit]
Description=Admiral Startup Configuration
Documentation=http://github.com/vmware/admiral
After=vic-appliance.target admiral-psc-token-ready.path
Requires=vic-appliance.target admiral-psc-token-ready.path

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /etc/vmware/admiral/configure_admiral.sh
ExecStartPost=/usr/bin/systemctl start get_token.service

[Install]
WantedBy=multi-user.target
