# installer/build/scripts/systemd/admiral/admiral_default_users.service
[Unit]
Description=Admiral Default Users
Documentation=http://github.com/vmware/admiral
After=vic-appliance-ready.target admiral.service get_token.service
Requires=vic-appliance-ready.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /etc/vmware/admiral/add_default_users.sh

[Install]
WantedBy=multi-user.target
