# installer/build/scripts/systemd/admiral/admiral_startup.path
[Unit]
Description=Admiral PSC token dependency
Documentation=https://github.com/vmware/vic-product/installer

[Path]
PathExists=/etc/vmware/psc/admiral/admiral_startup

[Install]
WantedBy=paths.target
