# installer/build/scripts/systemd/admiral/admiral-psc-token-ready.path
[Unit]
Description=Admiral PSC token dependency
Documentation=https://github.com/vmware/vic-product/installer

[Path]
PathExists=/etc/vmware/psc/admiral/psc-config.properties
Unit=admiral.service

[Install]
WantedBy=paths.target
