# installer/build/scripts/systemd/mounts/vic-mounts.target
[Unit]
Description=Disk mounts required by VIC(Integrated containers engine on vSphere)
Documentation=https://github.com/vmware/vic/blob/master/doc/bundle/OVA-spec.md
Requires=storage-data.mount storage-db.mount storage-log.mount
After=storage-data.mount storage-db.mount storage-log.mount

[Install]
WantedBy=vic-appliance.target
