# installer/build/scripts/systemd/harbor/harbor-psc-token-ready.path
[Unit]
Description=Harbor PSC token dependency
Documentation=https://github.com/vmware/vic-product/installer

[Path]
PathExists=/etc/vmware/psc/harbor/harbor_startup
Unit=harbor.service

[Install]
WantedBy=paths.target
