# installer/build/scripts/fileserver/fileserver.service
[Unit]
Description=VIC Unified Installer Web Server
Documentation=https://github.com/vmware/vic
After=vic-appliance.target
Requires=vic-appliance.target

[Service]
Type=simple
Restart=on-failure
RestartSec=15
ExecStartPre=/usr/bin/bash /etc/vmware/fileserver/configure_fileserver.sh
ExecStart=/etc/vmware/fileserver/start_fileserver.sh

[Install]
WantedBy=vic-appliance.target